# homelab/01-server-init/files/homelab-boot-notify.service
[Unit]
Description=Homelab Boot Notification Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelab-boot-notify.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target